# piclaw/piclaw-src/systemd/piclaw-watchdog.service
[Unit]
Description=PiClaw Watchdog – System Security & Health Monitor
# Starts independently – does NOT depend on piclaw-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# !! SEPARATE USER !! – piclaw user cannot write watchdog logs
User=piclaw-watchdog
Group=piclaw-watchdog
WorkingDirectory=/var/lib/piclaw-watchdog

ExecStart=/opt/piclaw/.venv/bin/python -m piclaw.agents.watchdog

Restart=always
RestartSec=15
StartLimitBurst=10
StartLimitIntervalSec=300

StandardOutput=append:/var/log/piclaw-watchdog/watchdog.log
StandardError=append:/var/log/piclaw-watchdog/watchdog.log

Environment="PICLAW_CONFIG=/etc/piclaw/config.toml"

# Hardened: watchdog cannot be stopped by piclaw user
ProtectSystem=strict
ReadWritePaths=/etc/piclaw/ipc /var/log/piclaw-watchdog /var/lib/piclaw-watchdog
ReadOnlyPaths=/etc/piclaw
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
