# piclaw/piclaw-src/systemd/piclaw-api.service
[Unit]
Description=PiClaw Web API & Dashboard (port 7842)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=piclaw
Group=piclaw
WorkingDirectory=/opt/piclaw

ExecStart=/opt/piclaw/.venv/bin/python -m uvicorn piclaw.api:app \
    --host 0.0.0.0 \
    --port 7842 \
    --log-level info \
    --no-access-log

Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

StandardOutput=append:/var/log/piclaw/api.log
StandardError=append:/var/log/piclaw/api.log

Environment="PICLAW_CONFIG=/etc/piclaw/config.toml"
Environment="PYTHONPATH=/opt/piclaw"

PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
