# piclaw/piclaw-src/systemd/piclaw-agent.service
[Unit]
Description=PiClaw AI Agent Daemon
After=network-online.target piclaw-api.service
Wants=network-online.target

[Service]
Type=simple
User=piclaw
Group=piclaw
WorkingDirectory=/home/piclaw

ExecStart=/usr/local/bin/piclaw-agent
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

StandardOutput=append:/var/log/piclaw/agent.log
StandardError=append:/var/log/piclaw/agent.log

Environment="PICLAW_CONFIG=/etc/piclaw/config.toml"

PrivateTmp=true
NoNewPrivileges=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
