# piclaw-os/systemd/piclaw-crawler.service
[Unit]
Description=PiClaw WebCrawler Sub-Agent
After=network-online.target piclaw-api.service
Wants=network-online.target

[Service]
Type=simple
# Same user as mainagent – orchestrated via jobs.db
User=piclaw
Group=piclaw
WorkingDirectory=/opt/piclaw

ExecStart=/opt/piclaw/.venv/bin/python -m piclaw.agents.crawler

Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=120

StandardOutput=append:/var/log/piclaw/crawler.log
StandardError=append:/var/log/piclaw/crawler.log

Environment="PICLAW_CONFIG=/etc/piclaw/config.toml"

PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
